# Linux Kernel/USB PlugIn/10-mert-kaplan.rules
ACTION=="add",KERNEL=="sdb*", RUN+="/home/debian/autorun.sh"